# global/overlay/etc/cron.d/cosmos
SHELL=/bin/sh
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

*/15 * * * *   root   test -f /etc/no-automatic-cosmos || (cosmos update && cosmos apply)
